# heroku/qbittorrent/services.d/syncthing/run
#!/bin/sh

if [ ! -z "$SYNCTHING" ]; then
	exec syncthing
else
        rm -rf /var/run/s6/services/syncthing
        s6-svscanctl -an /var/run/s6/services/
fi

