# heroku/citus/services.d/chiselclient/run
#!/usr/bin/with-contenv execlineb

chisel client $ROUTER $REMOTES