# heroku/qbittorrent/services.d/syncthing-inotify/run
#!/usr/bin/with-contenv sh

if [ ! -z "$SYNCTHING" ]; then
	exec syncthing-inotify
else
        rm -rf /var/run/s6/services/syncthing-inotify
        s6-svscanctl -an /var/run/s6/services/
fi
