# heroku/exec/services.d/exec/run
#!/usr/bin/with-contenv sh

eval "${EXEC}"