# heroku/exec/services.d/exec/run
#!/usr/bin/with-contenv sh

PATH=.:$PATH

if [ -n "${EXEC}" ]; then
    printf '%s' "${EXEC}" > start
    exec bash < start
else
    exec sleep infinity
fi
