# heroku/exec/services.d/exec/run
#!/usr/bin/with-contenv sh

PATH=.:$PATH
printf '%s' "${EXEC}" > start
exec bash start
