# heroku/chisel-alp/services.d/chisel/run
#!/usr/bin/with-contenv execlineb

chisel server
