# openshift/alpbase/services.d/dropbear/run
#!/usr/bin/with-contenv execlineb

multisubstitute {
        importas HOME HOME
        importas SSH_KEY SSH_KEY
}

foreground {
        mkdir -p ${HOME}/.ssh
}

foreground {
        redirfd -w 1 ${HOME}/.ssh/authorized_keys echo ${SSH_KEY}
}

dropbear -BRF -p 12322