# heroku/qbittorrent/services.d/qbittorrent/run
#!/usr/bin/with-contenv execlineb

/entrypoint.sh qbittorrent-nox

