# heroku/alpbase/services.d/dropbear/run
#!/usr/bin/with-contenv execlineb

redirfd -w 1 $HOME/.ssh/authorized_keys echo $SSH_KEY
    
dropbear -F